# systemd/pkgwatch-collect.service
[Unit]
Description=PkgWatch collect dpkg events into queue
# Prevent start-limit-hit during large installs
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=oneshot
ExecStart=/opt/pkgwatch/bin/pkgwatch-collect.sh
