# systemd/pkgwatch-audit.service
[Unit]
Description=PkgWatch auditd summary flush to Telegram

[Service]
Type=oneshot
ExecStart=/opt/pkgwatch/bin/pkgwatch-audit-flush.sh
