# systemd/pkgwatch-collect.path
[Unit]
Description=PkgWatch watch /var/log/dpkg.log

[Path]
PathChanged=/var/log/dpkg.log
Unit=pkgwatch-collect.service

[Install]
WantedBy=multi-user.target
