# systemd/pkgwatch-flush.timer
[Unit]
Description=Run PkgWatch flush periodically (rate-limited)

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=pkgwatch-flush.service

[Install]
WantedBy=timers.target
